# rpi_scripts/iptables_autosave/iptables-autosave.service
[Unit]
Description=Automatically save iptables rules

[Service]
Type=oneshot
ExecStart=/usr/sbin/netfilter-persistent save
